# cron/example.cron
* * * * * /usr/bin/php /var/www/eveotsv2/cron/bot.php >> /var/www/eveotsv2/cron/log/bot.log 2>&1
5 * * * * /usr/bin/php /var/www/eveotsv2/cron/checks/charactercheck.php >> /var/www/eveotsv2/cron/log/cron.log 2>&1
7 * * * * /usr/bin/php /var/www/eveotsv2/cron/checks/corporationcheck.php >> /var/www/eveotsv2/cron/log/cron.log 2>&1
9 * * * * /usr/bin/php /var/www/eveotsv2/cron/checks/alliancecheck.php >> /var/www/eveotsv2/cron/log/cron.log 2>&1
